# .deploy/systemd/vibe-orchestrator.service
[Unit]
Description=Vibe Orchestrator (Bridge-Agent)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/vibe/orchestrator.env
WorkingDirectory=/srv/vibe/bridge-agent
ExecStart=/usr/bin/env bash -lc 'exec npx next start -p "${PORT:-3002}" -H "${ORCH_BIND:-127.0.0.1}"'
Restart=on-failure
RestartSec=2
User=vibe
Group=vibe

[Install]
WantedBy=multi-user.target
