# .deploy/systemd/bridge-agent.service.d/override.conf
[Service]
# Load secrets and runtime env from /etc
EnvironmentFile=/etc/vibe/bridge-agent.env
Environment=NODE_ENV=production
